# files/etc/init.d/vpn-switch
#!/bin/sh /etc/rc.common

START=95
USE_PROCD=1

start_service() {
	/usr/bin/vpn-switch init
}

reload_service() {
	/usr/bin/vpn-switch reload
}
