# files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
# Disable serial login on the system console. MCU uses ttyS0.
#::askconsole:/usr/libexec/login.sh
